# ops/launchd/com.devloop.self-healing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.devloop.self-healing</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/somdax/dev/tiagosomda/dev-loop/src/backend/.venv/bin/python</string>
    <string>/Users/somdax/dev/tiagosomda/dev-loop/src/backend/devloop.py</string>
    <string>run</string>
    <string>self-healing</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/somdax/dev/tiagosomda/dev-loop</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Applications/ChatGPT.app/Contents/Resources</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>15</integer></dict>
  <key>StandardOutPath</key>
  <string>/Users/somdax/dev/tiagosomda/dev-loop/data/self-healing.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/somdax/dev/tiagosomda/dev-loop/data/self-healing.error.log</string>
</dict>
</plist>
